# deploy/something.service
[Unit]
Description=Something media platform (FastAPI)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/something
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/something/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
